# infrastructure/templates/shipsecure.service.j2
[Unit]
Description=ShipSecure Application
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory={{ app_directory }}
EnvironmentFile={{ app_directory }}/.env
ExecStart=/usr/bin/docker compose up -d
ExecStop=/usr/bin/docker compose down
ExecReload=/usr/bin/docker compose pull && /usr/bin/docker compose up -d
TimeoutStartSec=120
# Extended for graceful shutdown: Docker grace (90s) + buffer (5s) (Phase 24: INF-04)
TimeoutStopSec=95s
Restart=on-failure
RestartSec=15
User={{ deploy_user }}
Group={{ deploy_user }}

[Install]
WantedBy=multi-user.target
